# deploy-systemd/sipeed-cm5-fancontrol.service
[Unit]
Description=Sipeed CM5 Fan Controller
After=multi-user.target

[Service]
Type=simple
# Run as root to allow GPIO access
User=root
WorkingDirectory=/opt/sipeed-fancontrol
ExecStart=/usr/bin/python3 /opt/sipeed-fancontrol/fan_control.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
